# cron/autosave.cron
* * * * * ~/.config/tmux/autosave-sessions.zsh
